# docker/compose/rc.local
#!/bin/sh -e

# needed to run redis
echo never > /sys/kernel/mm/transparent_hugepage/enabled
sysctl vm.overcommit_memory=1
sysctl net.core.somaxconn=512

exit 0
